# nginx/nginx.container
# ~/.config/containers/systemd/webserver.container
[Unit]
Description=Nginx web server managed by Quadlet
After=network-online.target

[Container]
Image=docker.io/library/nginx:latest
ContainerName=webserver
PublishPort=8080:80

# Mount the website content
Volume=%h/website:/usr/share/nginx/html:ro,Z

# Health check
HealthCmd=curl -f http://localhost:80/ || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

# Auto-update when a new image is available
AutoUpdate=registry

[Service]
Restart=on-failure
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=default.target
